# olympics-tv.service
[Unit]
Description=Olympics TV FastAPI Application
After=network.target postgresql.service
Documentation=https://github.com/stosh99/olympics-tv

[Service]
Type=notify
User=olympics
Group=olympics
WorkingDirectory=/home/olympics/olympics-tv

# Load environment variables from .env
EnvironmentFile=/home/olympics/olympics-tv/.env

# Pre-start: ensure venv exists
ExecStartPre=/usr/bin/test -d /home/olympics/olympics-tv/venv

# Start the application with Gunicorn + Uvicorn
ExecStart=/home/olympics/olympics-tv/venv/bin/gunicorn api.main:app \
  -w 4 \
  -k uvicorn.workers.UvicornWorker \
  --bind 0.0.0.0:8000 \
  --access-logfile - \
  --error-logfile - \
  --log-level info

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
LimitNOFILE=65535
LimitNPROC=2048

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Standard output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=olympics-tv

[Install]
WantedBy=multi-user.target
